# rust-pi-app.service
[Unit]
Description=Rust Pi Zero Example Application
Documentation=https://github.com/your-username/rust_pi_zero_example
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/rust_pi_zero_example
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

# Network settings
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=10.0.0.0/8
IPAddressAllow=172.16.0.0/12
IPAddressAllow=192.168.0.0/16

[Install]
WantedBy=multi-user.target
